# vagrant-bootstrap/redis.conf
description "redis server"

start on runlevel [23]
stop on shutdown
expect daemon

respawn

env RUN_AS_USER=redis
env DAEMON=/usr/local/bin/redis-server
env CONFIGFILE=/etc/redis/redis-config.conf
env PIDFILE=/var/run/redis/redis.pid

pre-start script
  mkdir -p -m0766 /var/log/redis
  chown $RUN_AS_USER /var/log/redis
  mkdir -p -m0766 /var/run/redis
  chown $RUN_AS_USER /var/run/redis
  mkdir -p -m0766 /tmp/redis
  chown $RUN_AS_USER /tmp/redis
end script

script
  exec 2>>/dev/.initramfs/redis.log
  set -x 
  /sbin/start-stop-daemon --start --make-pidfile --pidfile $PIDFILE --chuid $RUN_AS_USER --exec $DAEMON $CONFIG_FILE >> /var/log/redis/redis.log 2>&1
end script