# vagrant-bootstrap/zookeeper.conf
# zookeeper server

description  "Zookeeper Server"

start on runlevel [2345]
stop on runlevel [!2345]
respawn
respawn limit 10 5
umask 022

kill timeout 30

env RUN_AS_USER=kafka-zookeeper
env BASEDIR=/usr/local/kafka
env PIDFILE=/var/run/zookeeper/zookeeper.pid
env DAEMON=/usr/local/kafka/bin/zookeeper-server-start.sh
env CONFIGFILE=/usr/local/kafka/config/zookeeper.properties

pre-start script
  test -x /usr/local/kafka/bin/zookeeper-server-start.sh || { stop; exit 0; }
  test -c /dev/null || { stop; exit 0; }
  mkdir -p -m0755 /var/log/zookeeper
  chown $RUN_AS_USER /var/log/zookeeper
  mkdir -p -m0755 /var/run/zookeeper
end script

script
    /sbin/start-stop-daemon --start --make-pidfile --pidfile $PIDFILE --chuid $RUN_AS_USER --exec $DAEMON $CONFIGFILE >> /var/log/zookeeper/zookeeper.log 2>&1
end script

pre-stop script
  /sbin/start-stop-daemon --stop --pidfile $PIDFILE --exec $DAEMON $DAEMON  >> /var/log/zookeeper/zookeeper.log 2>&1
end script